# templates/systemd/etc/systemd/system/sshnpd.service
[Unit]
Description=Ssh! No Ports Daemon
After=network-online.target

[Service]
User=<username>
Type=simple
Restart=always
ExecStartPre=/bin/sleep 10
RestartSec=3
ExecStart=/usr/bin/screen -DmS sshnpd %h/.local/bin/sshnpd -a <@client_atsign> -m <@device_manager_atsign> -u -d <device_name> -s -v

[Install]
WantedBy=multi-user.target