# packages/dart/sshnoports/bundles/shell/launchd/com.atsign.sshnpd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.atsign.sshnpd</string>
	<key>ProgramArguments</key>
    <array>
      <string>/Users/[user]/.local/bin/sshnpd</string>
      <string>-a</string>
      <string>[device_atsign]</string>
      <string>-m</string>
      <string>[client_atsign]</string>
      <string>-d</string>
      <string>[device_name]</string>
      <string>-su</string>
      <string>[policy_option]</string>
      <string>[policy_atsign]</string>
	</array>
  <key>KeepAlive</key>
  <true/>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
