# packages/sshnoports/bundles/shell/systemd/sshrvd.service
[Unit]
Description=Ssh No Ports Rendezvous Daemon
After=network-online.target

[Service]
User=<username>
Type=simple
Restart=always
RestartSec=3
ExecStart=/usr/local/bin/sshrvd -a <@atsign> -i <internet_address>

[Install]
WantedBy=multi-user.target